# deploy/netguard-runner.service
# NetGuard capture + scoring runner (privileged: needs CAP_NET_RAW).
#
# Install:
#   sudo cp deploy/netguard-runner.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now netguard-runner
#
# Prefer granting the capability to the venv interpreter rather than running as
# root:
#   sudo setcap cap_net_raw,cap_net_admin+eip \
#       $(readlink -f /opt/netguard/.venv/bin/python3)
[Unit]
Description=NetGuard capture and scoring runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netguard
Group=netguard
WorkingDirectory=/opt/netguard
Environment=NETGUARD_IFACE=eth0
# AmbientCapabilities grants CAP_NET_RAW without full root.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true
ExecStart=/opt/netguard/.venv/bin/python -m netguard.pipeline.runner --iface eth0
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
